# deploy/systemd/automaticcosmetic-backup.service
[Unit]
Description=AutomaticCosmetic runtime data backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/AutomaticCosmetic
EnvironmentFile=/root/AutomaticCosmetic/.env
ExecStart=/root/AutomaticCosmetic/.venv/bin/python scripts/backup_runtime_data.py --data-dir data --output-dir backups --env-path .env --retention-days 14
User=root

