# deploy/systemd/yclients-visit-confirmations.service
[Unit]
Description=AutomaticCosmetic daily visit confirmation sender
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/AutomaticCosmetic
Environment=PATH=/root/.local/bin:/root/AutomaticCosmetic/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/AutomaticCosmetic/.venv/bin/python /root/AutomaticCosmetic/scripts/send_visit_confirmations.py --quiet-empty

