# deploy/systemd/yclients-visit-confirmations.timer
[Unit]
Description=Run AutomaticCosmetic daily visit confirmation sender

[Timer]
OnCalendar=*-*-* 21:30:00
Persistent=true
Unit=yclients-visit-confirmations.service

[Install]
WantedBy=timers.target
