# deploy/systemd/yclients-avito-unanswered-monitor.service
[Unit]
Description=AutomaticCosmetic Avito unanswered monitor with delayed autoreply and Olga alerts
After=network-online.target yclients-avito-webhook.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/AutomaticCosmetic
Environment=PATH=/root/.local/bin:/root/AutomaticCosmetic/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/AutomaticCosmetic/.venv/bin/python /root/AutomaticCosmetic/scripts/avito_unanswered_monitor.py --autoreply --notify
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
