# deploy/systemd/yclients-yclients-integration.service
[Unit]
Description=Automatic Cosmetic YCLIENTS integration adapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/AutomaticCosmetic
Environment=PATH=/root/.local/bin:/root/AutomaticCosmetic/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/AutomaticCosmetic/run_yclients_integration.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
